# conf/initscript
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:          ufw
# Required-Start:    mountall.sh
# Required-Stop:     
# Default-Start:     S
# Default-Stop:      
# Short-Description: start firewall
### END INIT INFO

PATH="/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"

[ -x #PREFIX#/sbin/ufw ] || exit 0

. /lib/lsb/init-functions

if [ -s #CONFIG_PREFIX#/default/ufw ]; then
	. #CONFIG_PREFIX#/default/ufw
else
	log_failure_msg "Could not find #CONFIG_PREFIX#/default/ufw (disabling)"
	ENABLED="no"
fi

BEFORE_RULES="#CONFIG_PREFIX#/ufw/rules.before"
AFTER_RULES="#CONFIG_PREFIX#/ufw/rules.after"
IN_RULES="/var/lib/ufw/ufw-in.rules"

case "$1" in
start)
	if [ "$ENABLED" = "yes" ] || [ "$ENABLED" = "YES" ]; then
		if [ ! -z "$IPT_SYSCTL" ] && [ -s "$IPT_SYSCTL" ]; then
			log_action_begin_msg "ufw: Loading $IPT_SYSCTL"
			sysctl -q -p $IPT_SYSCTL
			log_action_end_msg $?
		fi
		for m in $IPT_MODULES
		do
			log_action_begin_msg "ufw: Loading $m"
			modprobe $m
			log_action_end_msg $?
		done
		log_action_begin_msg "ufw: Starting firewall"

		# flush the chains
		iptables -F

		# setup built-in chains' default policy
		iptables -P INPUT $DEFAULT_INPUT_POLICY
		iptables -P OUTPUT $DEFAULT_OUTPUT_POLICY
		iptables -P FORWARD $DEFAULT_FORWARD_POLICY

		# setup some other chains that can be used later
		iptables -N ufw-not-local

		# setup ufw-before-* chains
		iptables -N ufw-before-input
		iptables -N ufw-before-output
		iptables -N ufw-before-forward
		iptables -A INPUT -j ufw-before-input
		iptables -A INPUT -j ufw-before-output
		iptables -A INPUT -j ufw-before-forward
		if [ -x "$RULES_BEFORE" ]; then
			iptables-restore < $RULES_BEFORE > /dev/null 2>&1
			if [ "$?" != "0" ]; then
				log_warning_msg "    Problem running 'rules.before'"
			fi
		fi
		iptables -A ufw-before-input -j RETURN
		iptables -A ufw-before-output -j RETURN
		iptables -A ufw-before-forward -j RETURN

		# setup ufw-user chain
		iptables -N ufw-user-input
		iptables -N ufw-user-output
		iptables -N ufw-user-forward
		iptables -A INPUT -j ufw-user-input
		iptables -A OUTPUT -j ufw-user-output
		iptables -A FORWARD -j ufw-user-forward
		if [ -s "$IN_RULES" ]; then
			iptables-restore < $IN_RULES > /dev/null 2>&1
		fi
		iptables -A ufw-user-input -j RETURN
		iptables -A ufw-user-output -j RETURN
		iptables -A ufw-user-forward -j RETURN
		
		# setup ufw-after-* chains
		iptables -N ufw-after-input
		iptables -N ufw-after-output
		iptables -N ufw-after-forward
		iptables -A INPUT -j ufw-after-input
		iptables -A INPUT -j ufw-after-output
		iptables -A INPUT -j ufw-after-forward
		if [ -x "$RULES_BEFORE" ]; then
			iptables-restore < $RULES_BEFORE > /dev/null 2>&1
			if [ "$?" != "0" ]; then
				log_warning_msg "    Problem running 'rules.after'"
			fi
		fi
		iptables -A ufw-after-input -j RETURN
		iptables -A ufw-after-output -j RETURN
		iptables -A ufw-after-forward -j RETURN

		log_action_end_msg $?
	else
		log_warning_msg "ufw: Skipping firewall (not enabled)"
	fi
	;;
stop)
	log_warning_msg "ufw: Please use 'ufw disable' to stop"
	;;
restart)
	if [ "$ENABLED" = "yes" ] || [ "$ENABLED" = "YES" ]; then
		log_action_begin_msg "ufw: Restarting firewall"
		#PREFIX#/sbin/ufw disable >/dev/null 2>/dev/null
		#PREFIX#/sbin/ufw enable >/dev/null 2>/dev/null
		log_action_end_msg $?
	else
		log_warning_msg "ufw: Skipping firewall (not enabled)"
	fi
	;;
force-reload)
	log_action_begin_msg "ufw: Forcibly re-loading firewall"
	#PREFIX#/sbin/ufw disable >/dev/null 2>/dev/null
	#PREFIX#/sbin/ufw enable >/dev/null 2>/dev/null
	log_action_end_msg $?
	;;
*)
	echo "Usage: /etc/init.d/ufw {start|stop|restart|force-reload}"
	exit 1
	;;
esac

exit 0

