# debian/postinst
#!/bin/sh -e

. /usr/share/debconf/confmodule

BEFORE_FILE="/usr/share/ufw/before.rules"
AFTER_FILE="/usr/share/ufw/after.rules"
case "$1" in
    configure)
	if [ -e "/etc/ufw/ufw.rules" ]; then
		db_get ufw/oldrules || true
	fi

        update-rc.d ufw start 39 S . > /dev/null

        if [ ! -e "$BEFORE_FILE" ] && [ -z "$2" ]; then
            # if no config and fresh install copy the template
            cp $BEFORE_FILE /etc/ufw
        fi

        if [ ! -e "$AFTER_FILE" ] && [ -z "$2" ]; then
            # if no config and fresh install copy the template
            cp $AFTER_FILE /etc/ufw
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#
