# debian/postinst
#!/bin/sh -e

case "$1" in
    configure)
        update-rc.d ufw start 39 S . > /dev/null
        if [ ! -e "/etc/ufw/ufw.rules" ]; then
            cp /usr/share/ufw/ufw.rules /etc/ufw
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#
