# debian/postinst
#!/bin/sh -e

. /usr/share/debconf/confmodule

RULES_PATH="/etc/ufw"
USER_PATH="/etc/ufw"
TEMPLATE_PATH="/usr/share/ufw"
case "$1" in
    configure)
        for f in before.rules before6.rules after.rules after6.rules
        do
	    # if migration or upgrade, copy over the files if they don't exist
            if [ -e "/etc/ufw/ufw.rules" ] || [ -z "$2" ]; then
                if [ ! -e "$RULES_PATH/$f" ] && [ -z "$2" ]; then
                    cp $TEMPLATE_PATH $RULES_PATH/$f
                fi
            fi
        done

        for f in user.rules user6.rules
        do
            if [ ! -e "$USER_PATH/$f" ]; then
                # if no config, copy the template
                cp $TEMPLATE_PATH $USER_PATH/$f
            fi
        done

        update-rc.d ufw start 39 S . > /dev/null
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#
