# etc/network/if-up.d/ufw
#!/bin/sh

[ "x$IFACE" != "xlo" ] && exit 0

[ -s /etc/default/ufw ] || exit 0
. /etc/default/ufw

if [ "$ENABLED" = "yes" ] || [ "$ENABLED" = "YES" ]; then
	/sbin/ufw enable
fi

